# rootfs/etc/services.d/dnsmasq/run
#!/usr/bin/env bash

if [ "$DEBUG" == "1" ]; then
    set -x

    whoami
fi

exec dnsmasq --conf-file=$DNSMASQ_CONF --keep-in-foreground
