# rootfs/etc/services.d/dnsmasq/run
#!/usr/bin/env bash

exec dnsmasq --conf-file=/etc/dnsmasq.d/dnsmasq.conf --keep-in-foreground
