# rootfs/etc/services.d/dnsmasq/run
#!/usr/bin/env bash

exec dnsmasq --conf-file=$DNSMASQ_CONF --keep-in-foreground
