# rootfs/etc/services.d/openvpn/run
#!/usr/bin/env bash

s6-setuidgid root

exec ovpn_run
