# rootfs/etc/services.d/openvpn/run
#!/usr/bin/execlineb -P

s6-setuidgid root

exec openvpn --mktun --dev tun

s6-setuidgid openvpn

exec ovpn_run
