# rootfs/etc/services.d/openvpn/run
#!/usr/bin/env bash

exec s6-setuidgid root openvpn --mktun --dev tun

exec ovpn_run
