# rootfs/etc/services.d/openvpn/run
#!/usr/bin/env bash

exec s6-setuidgid root ovpn_run
