# rootfs/etc/services.d/openvpn/run
#!/usr/bin/execlineb -P

s6-setuidgid root

exec ovpn_run
